# example/supervisor.conf
[program:echosrv]
# nginx switches user on its own
user=echosrv
directory=/home/echosrv/echoserver
command=bundle exec thin
redirect_stderr=true
autorestart=true
